# env/etc/udev/rules.d/fstab.rules
# Force fstab options for devices
KERNEL=="sd*", RUN+="/usr/sbin/rbfstab"
