# env/etc/udev/rules.d/fstab.rules
# Force fstab options for specific storage devices
RUN+="/usr/sbin/rbfstab"
